# systemd/pricedin-hourly.service
[Unit]
Description=Priced In hourly pass (options pricing + alerts)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/pricedin
EnvironmentFile=/opt/pricedin/.env
ExecStart=/opt/pricedin/.venv/bin/python -m workers.hourly
